# roles/library/files/booksonic.service
[Unit]
Description=Booksonic service

[Service]
WorkingDirectory=/opt/booksonic
ExecStart=/usr/bin/java -jar /opt/booksonic/booksonic.war
#There are a multitude of options here for setting memory usage, ports, etc. Using #this, the url you will visit is http://192.168.1.100:8080.
User=enrico
Group=enrico
Type=simple
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target