# roles/library/files/calibre-web.service
[Unit]
Description=Calibre-Web
#any descriptive name you like
After=remote-fs.target
#this prevents calibre from starting until after network shares are mounted. If you are storing everything locally, change remote to local 

[Service]
Type=simple
User=enrico
Group=enrico
ExecStart=/usr/bin/python3 /opt/calibre-web/cps.py
WorkingDirectory=/opt/calibre-web/

[Install]
WantedBy=multi-user.target
